# linux/raspberrypi/.config/systemd/user/frpc.service
[Unit]
Description = frpc
After=network.target
After=network-online.target

[Service]
ExecStart = /home/peng/TOOL/frp/frpc -c /home/peng/etc/frp/frpc.ini
ExecReload = /home/peng/TOOL/frp/frpc -c /home/peng/etc/frp/frpc.ini reload
Restart = always
RestartSec=5
Type = simple

[Install]
#WantedBy = multi-user.target
#WantedBy=network-online.target
WantedBy=default.target
